# crontab
AWS_PROFILE=pkoch-personal
#mh d m w
22 * * * * /home/pkoch/.pyenv/shims/aws s3 sync --exclude ".dropbox" --exclude "*/.dropbox" ~pkoch/Dropbox s3://pkoch-archive/Dropbox
